# daemon/x52d.service.in
[Unit]
Description=X52 driver daemon

[Service]
Type=simple
ExecStart=%bindir%/x52d -f -v
ExecReload=kill -HUP $MAINPID
